# scripts/etc/systemd/system/roulotte-syncgit.timer
[Unit]
Description=Timer quotidien - Sync Git Roulotte

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
